# peertube-server/templates/backup/peertube-backup.service.j2
[Unit]
Description=peertube-backup

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash {{ peertube.backups_directory }}/backup.sh