# peertube-server/templates/backup/peertube-backup.service
# TODO
[Unit]
Description=peertube-backup

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash TODO/backup.sh